# crontab
*/5 * * * * munin /etc/munin/TP-Link_TD-W9970-munin/getstats.sh
